# packages/RedHat/ntopd
#! /bin/sh
#
# Author: Pablo Ruiz <pruiz@ip6seguridad.com>
#
# chkconfig: 2345 53 23
# description: Network Monitor
#
# processname: ntopd
#

. /etc/rc.d/init.d/functions

if [ ! -f /etc/sysconfig/ntopd ]; then
    exit 0
fi

. /etc/sysconfig/ntopd

RETVAL=0
prog="ntopd"
PPATH="/usr/bin/ntop"

case "$1" in
    start)
 echo -n "Starting service ntop"
 daemon $PPATH  -i $NTOPD_IFACE -w $NTOPD_PORT $NTOP_FLAGS -P /var/lib/ntop
        RETVAL=$?
        [ "$RETVAL" = 0 ] && touch /var/lock/subsys/$prog
        echo
        ;;
    stop)
        echo -n $"Stopping $prog:"
        killproc $PPATH
        RETVAL=$?
        [ "$RETVAL" = 0 ] && rm -f /var/lock/subsys/$prog
        echo
        ;;
    restart|reload)
 $0 stop && $0 start
     ;;
    status)
        status $PPATH
        RETVAL=$?
        ;;
    condrestart)
        if [ -f /var/lock/subsys/$prog ] ; then
              $0  stop
              $0  start
         fi
         ;;
 *)
      echo "Usage: $0 {start|stop|restart|reload|status}"
     exit 1
  ;;
esac
exit $RETVAL


