# packages/debian/postrm
#!/bin/sh -e

set -e

if [ "$1" = "purge" ] ; then
  update-rc.d ntop remove >/dev/null
fi

exit 0
